# scripts/verynice.service.in
[Unit]
Description=A tool for dynamically adjusting the nice-level of processes

[Service]
ExecStart=@@BINDIR@@/verynice
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
